# troubleshooting/SiteSetupOnEc2/gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn.sock
SocketUser=script-runner

[Install]
WantedBy=sockets.target